# deploy/tgbot-ferzerkki.service
# Ручная установка (проще скрипт): из корня репозитория
#   sudo bash deploy/install-systemd.sh
#
# Ручной вариант: скопируй этот файл в /etc/systemd/system/ и поправь пути и User/Group.
# Логи: journalctl -u tgbot-ferzerkki -f

[Unit]
Description=TGBOTFerzerkki Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tgbot
Group=tgbot
WorkingDirectory=/opt/TGBOTFerzerkki
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/TGBOTFerzerkki/.venv/bin/python -m app.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
